# simulation/azure/flowproxy@.service
[Unit]
Description=FlowProxy orderflow proxy (%i)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/flowproxy/flowproxy.env
# Derive builder name/FQDN from the instance name
Environment=BUILDER_NAME=%i
Environment=BUILDERNET_NODE_NAME=%i
ExecStart=/usr/local/bin/flowproxy
Restart=always
RestartSec=3s
LimitNOFILE=65536

# Journald logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

